# systemd/artifact-stats.service
[Unit]
Description=LabDoctorM Artifact Stats — статистика и рейтинг артефактов
After=network.target

[Service]
Type=oneshot
ExecStartPre=/bin/mkdir -p /root/LabDoctorM/.qwen/artifacts
ExecStart=/bin/bash -c 'cd /root/LabDoctorM/projects/artifact-pulse && python3 artifact_stats.py --update 2>&1'
User=root
WorkingDirectory=/root/LabDoctorM
StandardOutput=journal
StandardError=journal
TimeoutSec=120

[Install]
WantedBy=multi-user.target
