# systemd/artifact-links.timer
[Unit]
Description=LabDoctorM Artifact Link Checker — каждые 12 часов

[Timer]
OnCalendar=*-*-* 00,12:15:00
Persistent=true
RandomizedDelaySec=300
FixedRandomDelay=true

[Install]
WantedBy=timers.target
