# systemd/artifact-aging.timer
[Unit]
Description=LabDoctorM Artifact Aging — ежедневно в 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
