# systemd/artifact-constraints.service
[Unit]
Description=LabDoctorM Artifact Constraints — проверка структурных ограничений
After=network.target

[Service]
Type=oneshot
ExecStartPre=/bin/mkdir -p /root/LabDoctorM/.qwen/artifacts
ExecStart=/bin/bash -c 'cd /root/LabDoctorM/projects/artifact-pulse && python3 artifact_constraints.py --json > /root/LabDoctorM/.qwen/artifacts/constraints_latest.json 2>&1'
User=root
WorkingDirectory=/root/LabDoctorM
StandardOutput=journal
StandardError=journal
TimeoutSec=120

[Install]
WantedBy=multi-user.target
