# systemd/artifact-provenance.service
[Unit]
Description=LabDoctorM Artifact Provenance — верификация и confidence decay
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'cd /root/LabDoctorM/projects/artifact-pulse && python3 src/artifact_provenance.py --refresh 2>&1'
User=root
WorkingDirectory=/root/LabDoctorM
StandardOutput=journal
StandardError=journal
TimeoutSec=120

[Install]
WantedBy=multi-user.target
