# systemd/artifact-aging.service
[Unit]
Description=LabDoctorM Artifact Aging ‚ÄĒ —Ā—ā–į—Ä–Ķ–Ĺ–ł–Ķ –į—Ä—ā–Ķ—Ą–į–ļ—ā–ĺ–≤ (active‚Üístale‚Üíarchived)
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'cd /root/LabDoctorM/projects/artifact-pulse && python3 src/artifact_aging.py --dry-run 2>&1'
User=root
WorkingDirectory=/root/LabDoctorM
StandardOutput=journal
StandardError=journal
TimeoutSec=120

[Install]
WantedBy=multi-user.target
