# systemd/artifact-links.timer
[Unit]
Description=LabDoctorM Artifact Link Checker — каждые 12 часов

[Timer]
OnCalendar=*-*-* 00,12:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
