# systemd/search-rebuild.timer
[Unit]
Description=LabDoctorM Search Index Rebuild — каждые 6 часов

[Timer]
OnCalendar=*-*-* 02,08,14,20:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
